# scripts/feedemy-printer.service
[Unit]
Description=Feedemy Printer Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/feedemy-printer
ExecStart=/opt/feedemy-printer/venv/bin/python -m src.main
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
